# docker/openvswitch/openvswitch-ipsec/s6-rc.d/ipsec_check_nflog/up
#!/command/execlineb -P
/usr/sbin/ipsec --checknflog
