# docker/openvswitch/openvswitch-ipsec/s6-rc.d/ipsec_check_modules/up
#!/command/execlineb -P
/usr/libexec/ipsec/_stackmanager start
