# docker/openvswitch/openvswitch-ipsec/s6-rc.d/ike_daemon/run
#!/command/execlineb -P
/usr/libexec/ipsec/pluto \
    --leak-detective \
    --config /etc/ipsec.conf \
    --nofork \
    --nssdir /etc/ipsec.d \
    --logfile /var/log/kolla/openvswitch/pluto.log
