# docker/openvswitch/openvswitch-ipsec/s6-rc.d/ipsec_monitor/run
#!/command/execlineb -P
/usr/share/openvswitch/scripts/ovs-monitor-ipsec \
    --no-restart-ike-daemon \
    --ike-daemon=libreswan \
    --pidfile=/run/openvswitch/ovs-monitor-ipsec.pid \
    --ipsec-ctl=/run/openvswitch/pluto.ctl \
    --log-file=/var/log/kolla/openvswitch/ovs-monitor-ipsec.log
    unix:/run/openvswitch/db.sock
