# docker/openvswitch/openvswitch-ipsec/s6-rc.d/ipsec_monitor/run
#!/command/execlineb -P
/usr/share/openvswitch/scripts/ovs-ctl \
    --no-restart-ike-daemon \
    --ike-daemon=libreswan \
    start-ovs-ipsec
