# docker/openvswitch/openvswitch-ipsec/s6-rc.d/ipsec_check_config/up
#!/command/execlineb -P
/usr/libexec/ipsec/addconn \
    --config /etc/ipsec.conf \
    --checkconfig
