# docker/openvswitch/openvswitch-ipsec/s6-rc.d/kolla_start/up
#!/command/execlineb -S0
/command/with-contenv kolla_start
